# supervisord.conf
[supervisord]
nodaemon=true

[program:cats]
command=python3 ./cats/app.py

[program:api]
command=python3 ./flask-mongo/hello.py
